# services/hubot/run
#!/usr/bin/env bash

cd /opt/myhubot
/opt/myhubot/bin/hubot -a slack -n hubot
