# services/hubot/run
#!/usr/bin/env bash

${HUBOT_PATH}/bin/hubot -a slack -n ${HUBOT_NAME}
