# services/hubot/run
#!/usr/bin/env bash

exec /sbin/setuser ${HUBOT_USER} ${HUBOT_PATH}/bin/hubot -a slack -n ${HUBOT_NAME}
